# counter-app.service
[Unit]
Description=Counter Tracker — real-time counter tracking app
Documentation=https://github.com/yourusername/counter-tracker
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/counter-tracker/backend
Environment="PATH=/home/pi/counter-tracker/backend/venv/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/pi/counter-tracker/backend/venv/bin/python3 main.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=counter-tracker

# Resource limits
CPUQuota=50%

# Basic hardening
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
